# packages/notecoin/notecoin-0.11.24-py3.9.egg/notecoin/server/server.ini
[program:notecoin_server]
directory = .
command = python notecoin_server.py 
autostart = true
autorestart = true
user = bingtao
stdout_logfile = /notechats/logs/notecoin/strategy.log
